# passwd.service
[Unit]
Description=Netwerkfix-sysprep

[Service]
#ExecStartPre=
ExecStart=/etc/sysprep/password.sh
SyslogIdentifier=Diskutilization
#ExecStop=

[Install]
WantedBy=multi-user.target
